# config/s6/workspace-fix/up
#!/bin/bash
# Setup workspace and runtime directory

echo "Setting up runtime directory..."
mkdir -p /tmp/runtime-jovian
chown jovian:jovian /tmp/runtime-jovian
chmod 700 /tmp/runtime-jovian
echo "Runtime directory setup complete"