# config/s6/workspace-fix/up
#!/command/execlineb -P
foreground { echo "Setting up runtime directory and permissions..." }
foreground { mkdir -p /tmp/runtime-jovian }
foreground { chown -R jovian:jovian /tmp/runtime-jovian /workspace /home/jovian }
chmod 700 /tmp/runtime-jovian